# tomcat.conf
[program:tomcat]
command=/opt/tomcat-latest/bin/catalina.sh run
startsecs=5
stopsignal=INT
user=tomcat9
redirect_stderr=true
stdout_logfile=/var/log/tomcat.log

